# dist/etc/skel/.config/autostart/Picom.desktop
[Desktop Entry]
Exec=picom
X-XFCE-Autostart-Override=true
Hidden=false
